# api/cron/crontab
*/30 * * * * python3 /code/manage.py test
