# api/cron/crontab
*/30 * * * * . /root/project_env.sh;python   /code/manage.py test
